# Kvantkode/scripts/launchd/com.void.update-python.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.void.update-python-weekly</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>node</string>
      <string>/Users/yashasnaidu/Kvantcode/void/scripts/update-python-weekly.js</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PY_TRADING_PATH</key><string>/Users/yashasnaidu/Kvantcode/python</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>3</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key><string>/tmp/void-python-weekly.out</string>
    <key>StandardErrorPath</key><string>/tmp/void-python-weekly.err</string>
    <key>RunAtLoad</key><true/>
  </dict>
</plist>
